# system-daemon/softu2f.socket
[Socket]
ListenStream=/run/softu2f/softu2f.sock

[Install]
WantedBy=sockets.target